# pi-zero/systemd/dhcpcd.service
[Unit]
Description=DHCP Client Daemon
After=network.target
Wants=network.target

[Service]
Type=forking
PIDFile=/run/dhcpcd.pid
ExecStart=/usr/sbin/dhcpcd -w
ExecStop=/usr/sbin/dhcpcd -x
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 