# systemd/offgridnet.service
[Unit]
Description=OffGridNet Flask Backend
After=network.target hostapd.service dnsmasq.service
Requires=network.target
Wants=hostapd.service dnsmasq.service
Conflicts=nginx.service apache2.service

[Service]
User=root
Group=root
WorkingDirectory=/home/sunny/offgridnet-node
Environment=PATH=/home/sunny/offgridnet-node/backend/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/sunny/offgridnet-node/backend/venv/bin/python backend/app.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/var/log/offgridnet/flask.log
StandardError=append:/var/log/offgridnet/flask-error.log
Environment=FLASK_APP=backend/app.py
Environment=FLASK_ENV=production
Environment=FLASK_PORT=80

[Install]
WantedBy=multi-user.target 