# systemd/set-wlan-ip.service
[Unit]
Description=Configure wireless interface and assign static IP for OffGridNet
Before=hostapd.service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Stop any conflicting services
ExecStart=/bin/bash -c 'systemctl stop wpa_supplicant'
ExecStart=/bin/bash -c 'systemctl disable wpa_supplicant'
# Unblock WiFi and configure interface
ExecStart=/bin/bash -c 'rfkill unblock wifi'
ExecStart=/bin/bash -c 'ip link set wlan0 down'
ExecStart=/bin/bash -c 'sleep 2'
ExecStart=/bin/bash -c 'iw dev wlan0 set type __ap || true'
ExecStart=/bin/bash -c 'sleep 2'
ExecStart=/bin/bash -c 'ip link set wlan0 up'
ExecStart=/bin/bash -c 'iw reg set US'
ExecStart=/bin/bash -c 'iw wlan0 set power_save off'
# Wait and assign IP
ExecStart=/bin/bash -c "/bin/sleep 10 && /sbin/ip addr add 192.168.4.1/24 dev wlan0"
RemainAfterExit=true

[Install]
WantedBy=multi-user.target 